# tailwind.service
# Better: use envs, create separate home directory, user and enable hardening
# Current config uses root for quick iteration

[Unit]
Description=Tailwind API Service
After=network.target
Wants=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/root/project-tailwind
EnvironmentFile=/etc/tailwind.env

# Security hardening (switch to full/true for production)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=off
ProtectHome=false

# Run the script with the env's Python
ExecStart=/root/miniconda3/envs/sdrizzle/bin/python /root/project-tailwind/run_server.py
Restart=never
RestartSec=3
TimeoutStopSec=30
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
